# hotlapping/crontab
1 5 * * * ~/restart_server.sh >> ~/cron.log 2>&1